# vl-ocr.service
[Unit]
Description=PaddleOCR-VL-1.5 vLLM Server
After=network.target

[Service]
User=pedro
Environment=HSA_OVERRIDE_GFX_VERSION=11.0.0
Environment=PYTORCH_HIP_ALLOC_CONF=expandable_segments:True
Environment=VLLM_TARGET_DEVICE=rocm
Environment=ROCM_PATH=/opt/rocm
ExecStart=/home/pedro/Lobeira/.venv/bin/vllm serve PaddlePaddle/PaddleOCR-VL-1.5 \
  --dtype float16 \
  --device rocm \
  --gpu-memory-utilization 0.85 \
  --max-num-batched-tokens 8192 \
  --port 8100 \
  --host 127.0.0.1
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
